# com.things-api.server.plist
<!--
  Things API launchd agent

  Install:
    1. Edit this file — set WorkingDirectory and ProgramArguments paths
    2. Ensure your .env file exists in the working directory with THINGS_API_TOKEN set
    3. cp com.things-api.server.plist ~/Library/LaunchAgents/
    4. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.things-api.server.plist

  Uninstall:
    launchctl bootout gui/$(id -u)/com.things-api.server
    rm ~/Library/LaunchAgents/com.things-api.server.plist
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.things-api.server</string>

	<!-- Working directory — the server reads .env from here -->
	<key>WorkingDirectory</key>
	<string>/path/to/things-api</string>

	<!--
	  Option A: Run from local source (development / pre-PyPI)
	  Uses uv run to execute from the project directory.
	  Replace /path/to/uv with the output of: which uv
	-->
	<key>ProgramArguments</key>
	<array>
		<string>/path/to/uv</string>
		<string>run</string>
		<string>things-api</string>
	</array>

	<!--
	  Option B: Run from PyPI (after publishing)
	  Uncomment this block and comment out Option A above.
	  Replace /path/to/uvx with the output of: which uvx

	<key>ProgramArguments</key>
	<array>
		<string>/path/to/uvx</string>
		<string>things-api</string>
	</array>
	-->

	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/things-api.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/things-api.log</string>
</dict>
</plist>
