# com.things-api.server.plist
<!--
  Things API launchd agent

  Prerequisites:
    1. Run 'things-api init' first to create your config at ~/.config/things-api/config
    2. Replace /path/to/uvx below with the output of: which uvx

  Install:
    cp com.things-api.server.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.things-api.server.plist

  Uninstall:
    launchctl bootout gui/$(id -u)/com.things-api.server
    rm ~/Library/LaunchAgents/com.things-api.server.plist
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.things-api.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/path/to/uvx</string>
		<string>things-api</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/things-api.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/things-api.log</string>
</dict>
</plist>
